# templates/default/upstart.fastcgi.conf.erb
description "php-fastcgi - respawning UNIX Socket"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

expect fork
respawn

env PIDFILE=/var/run/php-fastcgi.pid
env CHILDREN=2
env RUNNER=/usr/local/bin/fcgiwrap

exec /usr/bin/spawn-fcgi -s <%= node[:nginx][:fcgi_socket] %> -P $PIDFILE -C $CHILDREN -f $RUNNER -u <%= node[:nginx][:user] %> -g <%= node[:nginx][:group] %>
